# pkg/uptimerec/usr/lib/systemd/system/uptimerec_cal.service
[Unit]
Description=record uptime cal

[Service]
Type=simple
ExecStart=/usr/bin/uptimerec.py

[Install]
WantedBy=multi-user.target
