# uptimerec/uptimerec.service
[Unit]
Description=record uptime

[Service]
Type=simple
Restart=always
ExecStart=/usr/bin/uptimerec.sh

[Install]
WantedBy=multi-user.target
